# examples/supervisor.conf
; Per-tenant supervisor program. Drop this into /etc/supervisor/conf.d/<tenant>-mail.conf
; and run `sudo supervisorctl reread && sudo supervisorctl update`.
;
; Replace the placeholders below.

[program:tenant-mail]
command=python3 /opt/catchall-mail-stack/app.py
directory=/opt/catchall-mail-stack
user=mailuser
autostart=true
autorestart=true
stderr_logfile=/var/log/tenant-mail.err.log
stdout_logfile=/var/log/tenant-mail.out.log
environment=
    MAIL_PORT="8517",
    SMTP_PORT="2525",
    SMTP_HOST="0.0.0.0",
    MAIL_DOMAIN="example.com",
    MAIL_HOSTNAME="mail.example.com",
    MAIL_BRAND="Example Mail",
    MAIL_COOKIE="example_mail_auth",
    MAIL_DB_DSN="postgresql://mailuser@/examplemail",
    MAIL_FORWARD_TO="you@gmail.com",
    MAIL_DEFAULT_FROM="info@example.com",
    MAIL_USER="admin",
    MAIL_PASS="CHANGE_ME",
    MAIL_SECRET="CHANGE_ME_TO_A_LONG_RANDOM_HEX_STRING",
    MAIL_ROOT_PATH="/mail",
    DKIM_KEY="/opt/catchall-mail-stack/keys/example.com.key",
    DKIM_SELECTOR="mail",
    HOME="/home/mailuser"

; --- Multi-tenant relay example -------------------------------------------------
; To run a SECOND tenant on the same VM (only one process can bind public :25):
;   1. Run a second supervisor program, identical to the above but with:
;        MAIL_PORT=8518, SMTP_PORT=2526, SMTP_HOST=127.0.0.1
;        MAIL_DOMAIN=other.com, etc., and a separate DB + DKIM key.
;   2. Add MAIL_RELAY_DOMAINS="other.com=127.0.0.1:2526" to the FIRST program's
;      environment so it relays inbound mail for other.com to the second app.
